# supervisord.conf
; Supervisord configuration for ObsidianPalace.
; Manages two processes:
;   1. obsidian-sync — continuous bidirectional vault sync via ob CLI
;   2. mcp-server   — FastAPI + MCP over SSE

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:obsidian-sync]
command=ob sync --continuous
directory=/data/vault
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-server]
command=uvicorn obsidian_palace.app:app --host 0.0.0.0 --port 8080
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
