# supervisord.conf
; ObsidianPalace process definitions for supervisord.
; Included by /etc/supervisor/supervisord.conf via conf.d/*.conf.
; Only [program:*] sections belong here — core supervisord config
; lives in the system-installed supervisord.conf.

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; nginx must start before it can proxy — but it will return 502 until
; uvicorn is up, which is fine (health check has start_period).
priority=100

[program:obsidian-sync]
command=/app/sync-guard.sh
directory=/data/vault
autostart=true
autorestart=true
startretries=2
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=200

[program:mcp-server]
command=uvicorn obsidian_palace.app:app --host 0.0.0.0 --port 8080
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=300
