# supervisord.conf
; Supervisord configuration for ObsidianPalace.
; Manages three processes:
;   1. nginx          — SSL termination + reverse proxy
;   2. obsidian-sync  — continuous bidirectional vault sync via ob CLI
;   3. mcp-server     — FastAPI + MCP over SSE

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; nginx must start before it can proxy — but it will return 502 until
; uvicorn is up, which is fine (health check has start_period).
priority=100

[program:obsidian-sync]
command=ob sync --continuous --path /data/vault
directory=/data/vault
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=200

[program:mcp-server]
command=uvicorn obsidian_palace.app:app --host 0.0.0.0 --port 8080
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=300
